# dnsmasq/supervisord.conf
[supervisord]
nodaemon=true

[program:dnsmasq]
command=/usr/sbin/dnsmasq --keep-in-foreground --log-facility=- --user=root --conf-dir=/etc/dnsmasq.d
